# roles/robokitty_devbox/templates/takopi.service.j2
[Unit]
Description=Takopi Telegram bridge for Robokitty Codex devbox
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ robokitty_bridge_user }}
Group={{ robokitty_bridge_user }}
WorkingDirectory={{ robokitty_workdir }}
Environment=HOME=/home/{{ robokitty_bridge_user }}
Environment=PATH={{ robokitty_wrappers_dir }}:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/bin:/home/{{ robokitty_bridge_user }}/.local/bin
ExecStart=/home/{{ robokitty_bridge_user }}/.local/bin/takopi
Restart=on-failure
RestartSec=5

# Takopi must cross to the Codex runner through the sudo wrapper.
# NoNewPrivileges would block sudo's required setuid transition.
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
{% if robokitty_enable_proc_hardening | bool %}
ProtectProc=invisible
{% endif %}
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
ReadWritePaths=/home/{{ robokitty_bridge_user }} {{ robokitty_workdir }} {{ robokitty_github_broker_exchange_dir }} {{ robokitty_worktree_source_dir }}
UMask=0077

[Install]
WantedBy=multi-user.target
